# src/crontab
*/5 * * * * /usr/bin/php /var/www/hostname/directory/email.php >/dev/null 2>&1